# vnc-supervisor.conf
[supervisord]
nodaemon=true

[program:xvfb]
command=/usr/bin/Xvfb :1 -screen 0 1920x1080x24
autorestart=true
priority=100
user=claude

[program:fluxbox]
command=/usr/bin/fluxbox
environment=DISPLAY=":1",HOME="/home/claude"
autorestart=true
priority=200
user=claude

[program:x11vnc]
command=/usr/bin/x11vnc -display :1 -forever -usepw -shared -rfbport 5900 -rfbauth /home/claude/.vnc/passwd
autorestart=true
priority=300
user=claude

[program:autocutsel]
command=/usr/bin/autocutsel -fork
environment=DISPLAY=":1"
autorestart=true
priority=250
user=claude

[program:firefox]
command=/usr/bin/firefox-esr --display=:1 --remote-debugging-port=9222 --remote-allow-origins=*
environment=DISPLAY=":1",MOZ_REMOTE_SETTINGS_DEVTOOLS="1",HOME="/home/claude"
autorestart=false
priority=400
user=claude

[program:terminal-api]
command=python3 /app/claude-terminal-controller.py
autorestart=true
priority=500
user=claude
environment=HOME="/home/claude"

[program:terminal]
command=/usr/bin/xterm -display :1 -geometry 120x30 -e /bin/bash
environment=DISPLAY=":1",HOME="/home/claude"
autorestart=false
priority=600
user=claude